# prod-server/etc/systemd/system/eterban.service
[Unit]
Description= Start ban.py

[Service]
ExecStart=/usr/share/eterban/ban.py

[Install]
WantedBy=multi-user.target