# gateway/etc/systemd/system/eterban.service
[Unit]
Description= Start eterban_switcher.py

[Service]
ExecStart=/etc/eterban/eterban_switcher.py

[Install]
WantedBy=multi-user.target